# openwrt/uconfig/files-render/etc/hotplug.d/ieee80211/20-uconfig-legacy
#!/bin/sh

if [ ! -f /lib/wifi/mac80211.uc ]; then
	ucode /lib/uconfig/wifi-detect.uc
fi
